# deploy/basketball-bot.service
[Unit]
Description=Basketball Bot Daemon (long-polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/opt/basketball-bot
ExecStart=/opt/basketball-bot/venv/bin/python bot_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения читает сам bot_daemon.py через python-dotenv (load_dotenv()).
# EnvironmentFile= здесь НЕ используется намеренно: systemd по-своему парсит кавычки
# и \n внутри значений, что ломает JSON в GOOGLE_SHEETS_CREDENTIALS (PEM-ключ
# получается повреждённым). python-dotenv парсит тот же файл корректно.

# Лимиты
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
